# config/launchd/com.yida.docserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for the per-host static doc server.

  Listens on 127.0.0.1:<host port> only; browsers reach it via the
  ssh -L tunnel that `sshs` already opens. See bin/docserver for the
  port/root resolution and config/sshs/hosts.json for the registry.

  Symlinked to ~/Library/LaunchAgents/com.yida.docserver.plist by
  install.sh. Manage with:
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.yida.docserver.plist
    launchctl bootout gui/$(id -u)/com.yida.docserver
    launchctl kickstart -k gui/$(id -u)/com.yida.docserver
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yida.docserver</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>exec "$HOME/.local/bin/docserver"</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/docserver.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/docserver.err.log</string>
</dict>
</plist>
